# crontab
MAILTO=tokuhirom@gmail.com

*/10 *  *   *   * perl /home/dotcloud/code/bin/watcher.pl

